# websrt.conf
[program:websrt]
command=/home/flibb/WebSRT/target/release/websrt-gateway --input srt --srt-mode listener --srt-port 9000 --bind 0.0.0.0 --latency 1000 --health-port 9090
directory=/home/flibb/WebSRT
autostart=true
autorestart=true
startretries=3
stderr_logfile=/home/flibb/WebSRT/logs/gateway.err.log
stdout_logfile=/home/flibb/WebSRT/logs/gateway.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
environment=RUST_LOG="debug"
